# s.conf
[program:s]
command=python3 /home/curator2/s/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/s/test.err.log
stdout_logfile=/home/curator2/s/test.out.log